# wake_on_lan/templates/wol.service.j2
# {{ ansible_managed }}

[Unit]
Description=Configure Wake On LAN

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ source_interface }} wol g

[Install]
WantedBy=basic.target
